# data/ubuntu-16.04-fresh-install-bash-hist.list
    1  sodo do-release-upgrade -d
    2  sudo do-release-upgrade -d
    3  ./stow-alle
    4  stow bash
    5  c
    6  exit
    7  c
    8  ilu
    9  sudo apt install stow synaptic
   10  sudp apt clean
   11  sudo apt clean
   12  sudo apt autoremove
   13  agi chromium-browser
   14  gsettings set org.gnome.desktop.interface show-battery-percentage true
   15  agi gir1.2-gconf-2.0
   16  iupdate
   17  update
   18  c
   19  upgradab;e
   20  cd $HOME/.icons
   21  mkdir .icons
   22  cd $HOME/.icons
   23  git clone https://github.com/keeferrourke/la-capitaine-icon-theme.git
   24  agi arc-theme
   25  sudo add-apt-repository ppa:tista/adapta
   26  sudo add-apt-repository ppa:libreoffice/ppa
   27  sudo add-apt-repository ppa:noobslab/themes
   28  timedatectl set-local-rtc 1
   29  sudo add-apt-repository -y ppa:stebbins/handbrake-releases
   30  sudo add-apt-repository ppa:varlesh-l/indicator-kdeconnect
   31  wget -q -O - https://dl-ssl.google.com/linux/linux_signing_key.pub | sudo apt-key add - 
   32  sudo sh -c 'echo "deb [arch=amd64] http://dl.google.com/linux/chrome/deb/ stable main" >> /etc/apt/sources.list.d/google-chrome.list'
   33  sudo apt-get install     apt-transport-https     ca-certificates     curl     software-properties-common
   34  curl -fsSL https://download.docker.com/linux/ubuntu/gpg | sudo apt-key add -
   35  udo add-apt-repository    "deb [arch=amd64] https://download.docker.com/linux/ubuntu \
   36     $(lsb_release -cs) \
   37  sudo add-apt-repository "deb [arch=amd64] https://download.docker.com/linux/ubuntu xenial stable
   38  sudo add-apt-repository "deb [arch=amd64] https://download.docker.com/linux/ubuntu xenial stable"
   39  agi docker-ce
   40  update
   41  sudo add-apt-repository ppa:webupd8team/indicator-kdeconnect
   42  update
   43  agi docker-ce
   44  sudo groupadd docker
   45  sudo usermod -aG docker $USER
   46  sudo systemctl enable docker
   47  c
   48  echo deb http://repository.spotify.com stable non-free | sudo tee /etc/apt/sources.list.d/spotify.list
   49  sudo apt-key adv --keyserver hkp://keyserver.ubuntu.com:80 --recv-keys BBEBDCB318AD50EC6865090613B00F1FD2C19886
   50  udo apt-key adv --keyserver keyserver.ubuntu.com --recv-keys ACCAF35C
   51  sudo apt-key adv --keyserver keyserver.ubuntu.com --recv-keys ACCAF35C
   52  echo deb ttp://apt.insynchq.com/[DISTRIBUTION] ^Con-free contrib | sudo tee /etc/apt/sources.list.d/spotify.list
   53  echo deb http://apt.insynchq.com/ubuntu zesty non-free contrib | sudo tee /etc/apt/sources.list.d/insync.list
   54  echo deb http://repository.spotify.com stable non-free | sudo tee /etc/apt/sources.list.d/spotify.list
   55  gi texlive-full apt-xapian-index ubuntu-restricted-extras ubuntu-restricted-extras nautilus-dropbox 7zip-full nautilus-admin playonlinux flashplugin-installer  variety icedtea-8-plugin openjdk-8-jre gdebi imwheel neofetch virtualbox virtualbox-ext-pack gnome-exe-thumbnailer gnome-raw-thumbnailer git gitg htop intel-microcode
   56  agi texlive-full apt-xapian-index ubuntu-restricted-extras ubuntu-restricted-extras nautilus-dropbox 7zip-full nautilus-admin playonlinux flashplugin-installer  variety icedtea-8-plugin openjdk-8-jre gdebi imwheel neofetch virtualbox virtualbox-ext-pack gnome-exe-thumbnailer gnome-raw-thumbnailer git gitg htop intel-microcode
   57  agi texlive-full apt-xapian-index ubuntu-restricted-extras ubuntu-restricted-extras nautilus-dropbox p7zip-full nautilus-admin playonlinux flashplugin-installer  variety icedtea-8-plugin openjdk-8-jre gdebi imwheel neofetch virtualbox virtualbox-ext-pack gnome-exe-thumbnailer gnome-raw-thumbnailer git gitg htop intel-microcode
   58  agi texlive-full apt-xapian-index ubuntu-restricted-extras ubuntu-restricted-extras nautilus-dropbox p7zip-full nautilus-admin flashplugin-installer  variety icedtea-8-plugin openjdk-8-jre gdebi imwheel neofetch virtualbox virtualbox-ext-pack gnome-exe-thumbnailer gnome-raw-thumbnailer git gitg htop intel-microcode
   59  c
   60  agi texlive-full apt-xapian-index ubuntu-restricted-extras ubuntu-restricted-extras nautilus-dropbox p7zip-full nautilus-admin flashplugin-installer variety icedtea-8-plugin openjdk-8-jre gdebi imwheel neofetch virtualbox virtualbox-ext-pack gnome-exe-thumbnailer gnome-raw-thumbnailer git gitg htop intel-microcode
   61  update
   62  agi apt-xapian-index
   63  update-apt-xapian-index -vf 
   64  sudo update-apt-xapian-index -vf 
   65  c
   66  e
   67  sudo add-apt-repository ppa:paulo-miguel-dias/mesa
   68  update
   69  agi audacity gimp wine-staging
   70  sudo dpkg --add-architecture i386 
   71  wget -nc https://dl.winehq.org/wine-builds/Release.key
   72  sudo apt-key add Release.key
   73  sudo apt-add-repository https://dl.winehq.org/wine-builds/ubuntu/
   74  update
   75  agi audacity gimp handbrake-gtk rawtherapee
   76  agi imwheel laptop-mode-tools pavucontrol gufw lm-sensors smartmontools
   77  agi adapta
   78  agi adapta-theme
   79  apt search adapta
   80  agi adapta-gtk-themes
   81  agi adapta-gtk-theme
   82  sudo apt-get install libc6:i386 libncurses5:i386 libstdc++6:i386 lib32z1 libbz2-1.0:i386
   83  sudo snap install atom --classic
   84  sudo snap remove atom --classic
   85  sudo snap remove atom 
   86  agi audacity gimp handbrake-gtk rawtherapee kdenlive
   87  agi audacity gimp handbrake-gtk rawtherapee kdenlive shutter
   88  sudo add-apt-repository ppa:papirus/papirus
   89  update
   90  papirus-icon-theme
   91  agi papirus-icon-theme
   92  ls
   93  sudo dpkg -i *.deb
   94  sudo apt-get install -f
   95  sudo apt install wine-staging playonlinux
   96  sudo apt install wine-staging 
   97  agi nautilus-terminal
   98  agi texstudio
   99  top
  100  neofetch
  101  c
  102  clinfo
  103  sudo nano /usr/share/applications/dropbox.desktop
  104  imwheel
  105  top
  106  sudo apt-get install gir1.2-gtop-2.0 gir1.2-networkmanager-1.0  gir1.2-clutter-1.0
  107  neofetch
  108  c
  109  neofetch
  110  c
  111  neofetch
  112  indicator-kdeconnect
  113  imwheel -h
  114  plist imwheel
  115  pslist imwheel
  116  alias
  117  pfilter imwheel
  118  chmod +r pixel-saver@deadalnix.me
  119  ls
  120  agi virtualbox
  121  agi virtualbox-ext-pack
  122  gnome-shell-extension-prefs
  123  XDG_CURRENT_DESKTOP=Unity
  124  dropbox start
  125  stow startup-items
  126  nano appindicator-fix-startup-dropbox.sh
  127  exit
  128  appindicator-fix-startup-dropbox
  129  xelatex
  130  latex
  131  agi xelatex
  132  agi texlive xelatex
  133  agi texlive-xelatex
  134  agi texlive-xetex
  135  agi texlive
  136  c
  137  indicator-kdeconnect
  138  sudo ufw allow 1714:1764/udp
  139  sudo ufw allow 1714:1764/tcp
  140  sudo ufw reload
  141  agi intel-ucode
  142  agi intel-microcode
  143  agi glxgears
  144  agi glx-gears
  145  sudo apt search  glx-gears
  146  sudo apt search  glx
  147  DRI_PRIME=1
  148  glxgears
  149  DRIPRIME=1
  150  glxgears
  151  DRI_PRIME=1 glxinfo | grep "OpenGL renderer"
  152  DRI_PRIME=1 glxgears | grep "OpenGL renderer"
  153  DRI_PRIME=1 glxgears
  154  clinfo
  155  agi darktable
  156  ./hashcat64.bin -b
  157  ./hashcat64.bin -b --force
  158  sudo ppa-purge paulo-miguel-dias/mesa
  159  top
  160  ./hashcat64.bin -b --force
  161  ./hashcat64.bin -b --force -d
  162  ./hashcat64.bin -b --force -d 1
  163  ./hashcat64.bin
  164  ./hashcat64.bin --help
  165  clinfo
  166  c
  167  clinfo
  168  c
  169  sudo add-apt-repository ppa:oibaf/graphics-drivers
  170  update
  171  upgrade
  172  apt upgrade
  173  sudo apt upgrade
  174  update
  175  sudo apt upgrade
  176  lmsensors
  177  sensors
  178  neofetch
  179  clinfo
  180  ./hashcat64.bin -b --force 
  181  top
  182  clinfo
  183  ./hashcat64.bin -b --force 
  184  top
  185  c
  186  neoftech --help
  187  neoftceh --help
  188  neofetch --help
  189  agi conky
  190  agi jq
  191  neofetch --help
  192  neofetch --options ubuntu
  193  neofetch --options 
  194  ls
  195  neofetch
  196  pfilter spotify
  197  killall spotify
  198  pfilter spotify
  199  restart gdm
  200  top
  201  pfilter dropbox
  202  killall dropbox
  203  dropbox
  204  dropbox -d
  205  dropbox start
  206  XDG_CURRENT_DESKTOP=Unity
  207  killall dropbox
  208  dropbox start
  209  pfilter kde
  210  rnm
  211  ls
  212  lsc
  213  c
  214  ls
  215  c
  216  ls
  217  cd
  218  ls
  219  sudo gedit /usr/share/applications/spotify.desktop
  220  c
  221  neofetch
  222  c
  223  neofetch
  224  top
  225  killall nylas
  226  sudo ./install
  227  agi matlab-support
  228  nano Text.txt
  229  touch Text.txt
  230  matlab -desktop
  231  alias
  232  dotfile app-launchers
  233  pfilter kde
  234  killall kde 
  235  killall kde* 
  236  killall *kde* 
  237  killall 
  238  killall kde connected
  239  killall kde connectd
  240  pfilter kde
  241  agi inkscape
  242  wine
  243  wine-staging
  244  wine -ver
  245  wine -version
  246  apt search wine
  247  winehq
  248  winehq-stable
  249  wine-staging
  250  wine-staging -ver
  251  wine-staging -version
  252  sudo apt install winbind
  253  snap
  254  snap find
  255  update
  256  upgradable
  257  apt autoremove
  258  sudo apt autoremove
  259  clinfo
  260  curl https://repo.skype.com/data/SKYPE-GPG-KEY | sudo apt-key add - 
  261  update
  262  apt cean
  263  apt clean
  264  sudo apt clean
  265  c
  266  sudo apt-get install libc6:i386 libncurses5:i386 libstdc++6:i386 lib32z1 libbz2-1.0:i386
  267  c
  268  history
  269  pfilter kde
  270  ./dfapt.sh
  271  c
  272  ./dfapt.sh
  273  c
  274  ./dfapt.sh
  275  c
  276  ./dfapt.sh
  277  c
  278  ./dfapt.sh
  279  c
  280  ./dfapt.sh
  281  c
  282  ./dfapt.sh
  283  c
  284  ./dfapt.sh
  285  c
  286  ./dfapt.sh
  287  c
  288  ./dfapt.sh
  289  c
  290  ./dfapt.sh
  291  c
  292  ./dfapt.sh
  293  c
  294  ./dfapt.sh
  295  c
  296  ./dfapt.sh
  297  c
  298  ./dfapt.sh
  299  c
  300  ./dfapt.shc
  301  ./dfapt.sh
  302  c
  303  ./dfapt.sh
  304  c
  305  ./dfapt.sh
  306  c
  307  ./dfapt.sh
  308  c
  309  ufw
  310  ufw status
  311  sudo ufw status
  312  history |grep ufw 
  313  v
  314  c
  315  ./dfapt.sh
  316  c
  317  ./dfapt.sh
  318  c
  319  ./dfapt.sh
  320  c
  321  ./dfapt.sh
  322  c
  323  ./dfapt.sh
  324  c
  325  ./dfapt.sh
  326  top
  327  dotfile app-launchers
  328  pfilter imwheel
  329  ls
  330  lspci -n | grep 0280
  331  modprobe iwn
  332  lspci | grep Wireless
  333  lspci | grep Wifi
  334  lspci | grep intel
  335  lspci | grep Centrino
  336  lspci -knn | grep Net -A2
  337  sudo tee /etc/modprobe.d/iwlwifi-opt.conf <<< "options iwlwifi 11n_disable=8"
  338  reboot
  339  gnome-shell-extension-tool
  340  gnome-shell-extension-tool --help
  341  gnome-shell-extension-tool -d sound-output-device-chooser@kgshank.net
  342  gnome-shell-extension-tool -e sound-output-device-chooser@kgshank.net
  343  c
  344  dotfile startup-items
  345  c
  346  fix-volume-extension-conflict 
  347  appindicator-fix-startup-dropbox 
  348  dropbox stop
  349  appindicator-fix-startup-dropbox 
  350  fix-volume-extension-conflict 
  351  pfilter kde
  352  fix-volume-extension-conflict 
  353  dropbox stop
  354  appindicator-fix-startup-dropbox 
  355  dropbox stop
  356  appindicator-fix-startup-dropbox 
  357  fix-volume-extension-conflict 
  358  appindicator-fix-startup-dropbox 
  359  c
  360  fix-volume-extension-conflict 
  361  neofetch
  362  c
  363  neofetch
  364  c
  365  neofetch
  366  fix-volume-extension-conflict 
  367  dotfile scripts
  368  fix-volume-extension
  369  pfilter wine
  370  pfilter exe
  371  fix-volume-extension
  372  gnome-shell-extension-tool -e sound-output-device-chooser@kgshank.net
  373  gnome-shell-extension-tool -r sound-output-device-chooser@kgshank.net
  374  sudo snap install canonical-livepatch
  375  sudo snap remove canonical-livepatch
  376  sudo update-desktop-database
  377  fix-volume-extension
  378  adb
  379  kvm-ok
  380  agi cpu-checker
  381  sudo apt autoremove
  382  sudo apt clean
  383  c
  384  update
  385  upgradable 
  386  sudo apt upgrade
  387  c
  388  sudo apt clean
  389  c
  390  adb
  391  c
  392  fix-volume-extension
  393  kvm-ok
  394  clinfo
  395  bluez
  396  blueman
  397  agi bueman
  398  agi blueman
  399  reboot
  400  apt-get install pulseaudio pulseaudio-module-bluetooth pavucontrol bluez-firmware
  401  sudo apt-get install pulseaudio pulseaudio-module-bluetooth pavucontrol bluez-firmware
  402  sudo apt-get install pulseaudio pulseaudio-module-bluetooth pavucontrol 
  403  gksudo gedit /etc/bluetooth/audio.conf
  404  sudo nano /etc/bluetooth/audio.conf 
  405  restart bluetooth
  406  pulseaudio -k
  407  sudo /etc/init.d/bluetooth restart
  408  fix-volume-extension
  409  history >> ~/Desktop/ubuntu-16.04-fresh-install-bash-hist.list
